# piweather/piweather-0.1.bin.files
/var/piweather/piweather.sqlite3
/etc/piweather/piweather.conf
/etc/init.d/piweather
/etc/default/piweather
/usr/bin/piweather
/usr/lib/libpiweather.so.1.0.0
/usr/lib/libpiweather.la
/usr/lib/libpiweather.a
/usr/include/piweather/sensors/SensorDevice.hpp
/usr/include/piweather/sensors/Humidity.hpp
/usr/include/piweather/sensors/Illuminance.hpp
/usr/include/piweather/sensors/ThresholdObservable.hpp
/usr/include/piweather/sensors/Barometer.hpp
/usr/include/piweather/sensors/Temperature.hpp
/usr/include/piweather/device/Connection.hpp
/usr/include/piweather/device/UidRegistry.hpp
/usr/include/piweather/device/IpAddress.hpp
/usr/include/piweather/device/Observer.hpp
/usr/include/piweather/device/Observable.hpp
/usr/include/piweather/device/Button.hpp
/usr/include/piweather/ip_connection.h
/usr/include/piweather/bricklet_lcd_20x4.h
/usr/include/piweather/bricklet_humidity.h
/usr/include/piweather/bricklet_ambient_light.h
/usr/include/piweather/bricklet_barometer.h
/usr/share/piweather/schema.sql
